# deploy/systemd/vocabume-worker-low.service
[Unit]
Description=VocabuMe low-priority Celery worker
After=network-online.target redis-server.service
Wants=network-online.target redis-server.service

[Service]
Type=simple
User=eduard
Group=eduard
WorkingDirectory=/home/eduard/englishbot
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/eduard/englishbot/venv/bin/celery -A core worker --loglevel=INFO --queues=vocabume-low --concurrency=1 --hostname=vocabume-low@%%h
Restart=on-failure
RestartSec=5
TimeoutStopSec=310

[Install]
WantedBy=multi-user.target
